# piper_boot_init.service
# Piper 机械臂开机初始化服务
#
# 安装方法：
#   sudo cp piper_boot_init.service /etc/systemd/system/
#   sudo cp build/piper_boot_init /opt/piper_tools/bin/
#   sudo systemctl daemon-reload
#   sudo systemctl enable piper_boot_init.service
#
# 查看日志：
#   journalctl -u piper_boot_init.service -f
#

[Unit]
Description=Piper Robot Arm Boot Initialization

# 等待 CAN 接口就绪
After=sys-subsystem-net-devices-can0.device
Requires=sys-subsystem-net-devices-can0.device

[Service]
Type=oneshot
RemainAfterExit=no

ExecStart=/opt/piper_tools/bin/piper_boot_init

Restart=on-failure
RestartSec=2s
StartLimitIntervalSec=30
StartLimitBurst=15

SyslogIdentifier=piper_boot_init

[Install]
WantedBy=multi-user.target
